# faces.service
[Unit]
Description=faces: распознавание лиц на Orange Pi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/orangepi/faces
ExecStart=/home/orangepi/faces/venv/bin/python main.py
Restart=on-failure
RestartSec=10
# Логи в файл (user journal на Orange Pi часто недоступен)
StandardOutput=append:/home/orangepi/faces/faces.log
StandardError=append:/home/orangepi/faces/faces.log

[Install]
WantedBy=default.target
